# scripts/scripts/systemd.timers/borgbackup.service
[Unit]
Description="Borg Backup Script"

[Service]
ExecStart=/home/backup/borg_backup.sh
