# scripts/scripts/systemd.timers/notes.service
[Unit]
Description="Check and update notes when needed"

[Service]
ExecStart=/home/throttlemeister/scripts/git-notes.sh
